# systemd/user/check-battery-user.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/bringsrain/.config/systemd/user/check-battery.sh
